# templates/centos/svscan.conf.erb
start on runlevel [12345]
respawn
exec <%= node['daemontools']['bin_dir'] %>/svscanboot
